# ops/systemd/stockaisle-backend.service
[Unit]
Description=StockAisle Backend API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ecom/backend
Environment=NODE_ENV=production
EnvironmentFile=/home/ubuntu/ecom/backend/.env
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
